# production/etc/systemd/system/csroom.service
[Unit]
Description=csroom backend service
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/csroom
Environment=PYTHONUNBUFFERED=1
Environment=USERS_JSON_FILE=/var/lib/csroom/users.json
ExecStartPre=+/bin/mkdir -p /var/lib/csroom/uploads /var/lib/csroom/classrooms
# Shallow chown only on the three top-level dirs. Do NOT recurse: user files
# inside are owned by the container user (999:995) and must stay that way so
# they remain writable from the terminal. Group is csroom-container (GID 995)
# so www-data (with supplementary 995) can read/write via group perms.
ExecStartPre=+/bin/chown www-data:csroom-container /var/lib/csroom /var/lib/csroom/uploads /var/lib/csroom/classrooms
# Setgid bit (2775) makes new entries inherit GID 995 automatically.
ExecStartPre=+/bin/chmod 2775 /var/lib/csroom /var/lib/csroom/uploads /var/lib/csroom/classrooms
ExecStart=/var/www/csroom/.venv/bin/python3 -m backend.api
Restart=always
RestartSec=5
User=www-data
Group=www-data
SupplementaryGroups=csroom-container
# CAP_CHOWN: chown of upload dirs to the container user.
# CAP_NET_ADMIN: insert iptables rules into DOCKER-USER to isolate per-user
# networks (backend/docker.py:90). Without it, iptables fails with
# "Permission denied (you must be root)" and host-egress isolation breaks.
AmbientCapabilities=CAP_CHOWN CAP_NET_ADMIN
# systemd's default soft fd limit is 1024. The backend opens pty + dockerd
# connection fds per terminal attach, and a leak in that path exhausted the
# limit in prod (EMFILE cascade → 500s on saves/reads/login). The leak is
# fixed in code but 1024 is too tight for a busy class anyway; 65536 gives
# comfortable headroom.
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
